# archiso/airootfs/etc/systemd/system/ly@tty1.service.d/override.conf
[Service]
# Set TTY color 0 (bg) and 7 (lightgrey/fg) to Dracula colors before Ly starts
ExecStartPre=/usr/bin/printf '%%b' '\e]P0282a36\e]P7f8f8f2\ec'
